# etc/supervisor.conf.in
[program:pypiproxy]
command = ${buildout:directory}/bin/gunicorn flask_pypi_proxy.views:app -b 0.0.0.0:8090 
process_name = pypiproxy
directory = ${buildout:directory}/bin
user=fifteen5
redirect_stderr=True
autostart=true
autorestart=true

